# backend_tasks/crontab
*/2 * * * * python /home/Austin/hackillinois-website/backend_tasks/main.py tasks